# miniarm/bluetoothtool/jasonbluetooth.service
[Unit]
Description=bluetooth supplicant

[Service]
ExecStart=/sbin/bluetoothstart

[Install]
WantedBy=multi-user.target

